# packaging/rhel/systemd/ciso-assistant-huey.service
[Unit]
Description=CISO Assistant Background Worker (Huey)
Documentation=https://github.com/intuitem/ciso-assistant-community
After=network.target ciso-assistant-backend.service
Requires=ciso-assistant-backend.service

[Service]
Type=simple
User=ciso-assistant
Group=ciso-assistant
WorkingDirectory=/opt/ciso-assistant/backend
EnvironmentFile=/etc/ciso-assistant/backend.env

# Use bundled Python virtualenv
ExecStart=/opt/ciso-assistant/venv/bin/python manage.py run_huey -w 2 -k process

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ciso-assistant/backend/db /opt/ciso-assistant/backend/media

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

[Install]
WantedBy=multi-user.target
